# docker-files/rootfs/etc/s6-overlay/s6-rc.d/phpfpm/finish
#!/bin/bash

# Exit the container if PHP FPM exits
kill -TERM 1
