# docker-files/rootfs/etc/s6-overlay/s6-rc.d/phpfpm/run
#!/command/execlineb

s6-envuidgid www-data
/usr/sbin/php-fpm8.4 --nodaemonize
