# docker-files/rootfs/etc/s6-overlay/s6-rc.d/mariadbd/finish
#!/bin/bash

# Exit the container if MariaDB exits
kill -TERM 1
