# docker-files/rootfs/etc/s6-overlay/s6-rc.d/mariadbd/run
#!/bin/sh

exec s6-setuidgid mysql mariadbd
