# docker-files/rootfs/etc/s6-overlay/s6-rc.d/apache2/finish
#!/bin/bash

# Exit the container if Apache exits
kill -TERM 1
